# home/private_dot_config/autostart/Keymapp.desktop
[Desktop Entry]
Type=Application
Name=Keymapp
Exec="/var/home/scta/.local/bin/keymapp"
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=10
