# home/private_dot_config/systemd/user/yubikey-touch-detector.socket
[Socket]
ListenStream=%t/yubikey-touch-detector.socket
RemoveOnStop=true

[Install]
WantedBy=sockets.target
